# mysql-stack.service
[Unit]
Description=MySQL Stack Deployment
After=network.target

[Service]
ExecStart=/var/lib/mysql/start.sh

[Install]
WantedBy=default.target